# services/xoul-discord.service
[Unit]
Description=Xoul Discord Bot Client
After=xoul.service
Wants=xoul.service

[Service]
Type=simple
WorkingDirectory=/root/xoul
ExecStart=/usr/bin/python3 /root/xoul/discord_client.py
Restart=always
RestartSec=10
Environment=XOUL_CONFIG=/root/xoul/config.json

[Install]
WantedBy=multi-user.target
